# cmd/diningUpdate/dining.timer
[Unit]
	Description=Check Dining

[Timer]
	OnUnitInactiveSec=30s
	Unit=dining.service

[Install]
	WantedBy=basic.target
