# cmd/diningSearch/dining.timer
[Unit]
	Description=Check Dining

[Timer]
#	OnUnitActiveSec=5min
	OnUnitInactiveSec=30s
	Unit=dining.service

[Install]
	WantedBy=basic.target
